# docs/agent-runs/live-product-infra-rpc/public-rpc-deployment-bundle/flowchain-supervisor.service.template
# FlowChain live service supervisor systemd template.
# Render on the owner host only. Keep rendered unit files and env files out of the repository.
[Unit]
Description=FlowChain live service supervisor
Wants=network-online.target
After=network-online.target flowchain-live.service

[Service]
Type=simple
WorkingDirectory=<FLOWCHAIN_REPO_ABSOLUTE_PATH>
User=<FLOWCHAIN_SERVICE_USER>
Group=<FLOWCHAIN_SERVICE_GROUP>
EnvironmentFile=<FLOWCHAIN_OWNER_ENV_FILE>
Environment=FLOWCHAIN_OWNER_ENV_FILE=<FLOWCHAIN_OWNER_ENV_FILE>
Environment=FLOWCHAIN_CONTROL_PLANE_CARGO_TARGET_DIR=<FLOWCHAIN_CONTROL_PLANE_CARGO_TARGET_DIR>
ExecStart=/usr/bin/env npm run flowchain:service:supervisor -- -IntervalSeconds 30 -MaxRestartAttempts 3
Restart=always
RestartSec=15
TimeoutStartSec=120
TimeoutStopSec=60
KillMode=process
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=<FLOWCHAIN_REPO_ABSOLUTE_PATH>/devnet <FLOWCHAIN_REPO_ABSOLUTE_PATH>/docs/agent-runs <FLOWCHAIN_REPO_ABSOLUTE_PATH>/services/bridge-relayer/out

[Install]
WantedBy=multi-user.target
