# docs/agent-runs/live-product-infra-rpc/public-rpc-deployment-bundle/flowchain-live.service.template
# FlowChain live service systemd template.
# Render on the owner host only. Keep rendered unit files and env files out of the repository.
[Unit]
Description=FlowChain live service
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=<FLOWCHAIN_REPO_ABSOLUTE_PATH>
User=<FLOWCHAIN_SERVICE_USER>
Group=<FLOWCHAIN_SERVICE_GROUP>
EnvironmentFile=<FLOWCHAIN_OWNER_ENV_FILE>
Environment=FLOWCHAIN_OWNER_ENV_FILE=<FLOWCHAIN_OWNER_ENV_FILE>
Environment=FLOWCHAIN_CONTROL_PLANE_CARGO_TARGET_DIR=<FLOWCHAIN_CONTROL_PLANE_CARGO_TARGET_DIR>
ExecStart=/usr/bin/env npm run flowchain:service:start -- -LiveProfile
ExecStartPost=/usr/bin/env npm run flowchain:service:status
ExecReload=/usr/bin/env npm run flowchain:service:restart -- -LiveProfile
ExecStop=/usr/bin/env npm run flowchain:service:stop
RemainAfterExit=yes
TimeoutStartSec=900
TimeoutStopSec=180
KillMode=process
Restart=on-failure
RestartSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=<FLOWCHAIN_REPO_ABSOLUTE_PATH>/devnet <FLOWCHAIN_REPO_ABSOLUTE_PATH>/docs/agent-runs <FLOWCHAIN_REPO_ABSOLUTE_PATH>/services/bridge-relayer/out

[Install]
WantedBy=multi-user.target
